# batglass.service
[Unit]
Description=BatGlass smart glasses daemon
# network.target = networking stack is up; intentionally not network-online.target
# so buttons work immediately on boot even without Wi-Fi.
# Live API requests should fail gracefully per button press when offline.
After=sound.target network.target
Wants=sound.target
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=simple
User=aloy
WorkingDirectory=/home/aloy/batglass

# Load API keys from the project .env (GEMINI_API_KEY, etc.)
EnvironmentFile=-/home/aloy/batglass/.env
Environment=PYTHONUNBUFFERED=1

# Kill any stale process from a previous unclean shutdown before starting.
ExecStartPre=-/usr/bin/pkill -u aloy -f "python -m buttons.__main__"
ExecStart=/home/aloy/batglass/run_buttons.sh

# Restart always so clean exits, crashes, and transient failures all recover.
Restart=always
RestartSec=10

StandardOutput=journal
StandardError=journal
SyslogIdentifier=batglass

[Install]
WantedBy=multi-user.target
